# systemd/previous-admin-backend.service
[Unit]
Description=Previous Admin Backend API
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=next
WorkingDirectory=/home/next/previous-admin
ExecStart=/usr/bin/npm run backend
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=previous-admin-backend

[Install]
WantedBy=multi-user.target
