# systemd/previous-admin-frontend.service
[Unit]
Description=Previous Admin Frontend
After=network.target previous-admin-backend.service
Wants=network-online.target

[Service]
Type=simple
User=next
WorkingDirectory=/home/next/previous-admin
ExecStart=/usr/bin/npm run preview -- --host 0.0.0.0 --port 2342
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=previous-admin-frontend

[Install]
WantedBy=multi-user.target
